# udev/90-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/bash /usr/local/bin/test.sh
