# udev/98-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/bash /opt/config/udev/video_out_switch.sh"
